# projects/Switch/packages/switch-gpu-profile/system.d/switch-gpu-profile.service
[Unit]
Description=Switch GPU Profile
ConditionKernelCommandLine=!installer
<<<<<<< HEAD
=======
Requires=xorg.service
After=xorg.service
>>>>>>> 654f30666a863476dfba1a85ce51a31089d23d6b

[Service]
Environment=TERM=linux
Type=oneshot
ExecStart=/usr/bin/gpu-profile init
StandardError=syslog
RemainAfterExit=no

[Install]
WantedBy=graphical.target
