# packages/lakka/Switch/system.d/tmp-assets.mount
[Unit]
Description=Assets directory
<<<<<<< HEAD
Before=retroarch.target
=======
Before=retroarch.service
>>>>>>> 654f30666a863476dfba1a85ce51a31089d23d6b
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=/storage/assets:/usr/share/retroarch-assets
Where=/tmp/assets
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
