# packages/lakka/Switch/system.d/tmp-shaders.mount
[Unit]
Description=Shaders directory RetroArch
<<<<<<< HEAD
Before=retroarch.target
=======
Before=retroarch.service
>>>>>>> 654f30666a863476dfba1a85ce51a31089d23d6b
After=storage.mount
After=systemd-tmpfiles-setup.service

[Mount]
What=/storage/shaders:/usr/share/common-shaders
Where=/tmp/shaders
Type=mergerfs
Options=defaults,allow_other,use_ino

[Install]
WantedBy=retroarch.target
